# templates/teamcity-agent.erb
#! /bin/sh
### BEGIN INIT INFO
# Provides: <%= @service %>
# Required-Start: $local_fs $remote_fs $syslog $network
# Required-Stop: $local_fs $remote_fs $syslog $network
# Default-Start: 2 3 4 5
# Default-Stop: 0 1 6
# Short-Description: TeamCity Build Agent
# Description: TeamCity Build Agent
### END INIT INFO

set -e

PATH=/sbin:/bin:/usr/bin
USER=<%= @user %>
SCRIPT=<%= @home %>/bin/agent.sh
NAME=<%= @agent_name %>
PIDFILE=<%= @temp_dir %>/buildAgent.pid

start() {
  if [ status ];
  then
    echo "teamcity-agent is already running"
  else
    echo "starting teamcity-agent"
    rm -f "$PIDFILE"
    su -s /bin/bash -c "$SCRIPT start" $USER
  fi
}

stop() {
  echo "stopping teamcity-agent"
  su -s /bin/bash -c "$SCRIPT stop" $USER
}

status() {
  if [ -e "$PIDFILE" ];
  then
    ps -p $(cat $PIDFILE) > /dev/null
    RUNNING=$?
  else
    RUNNING=1
  fi
  if [ $RUNNING == 0 ]; 
  then 
    echo "teamcity-agent is running"; 
  fi
  return $RUNNING
}

case "$1" in
  start)
    start
    ;;
  stop)
    stop
    ;;
  restart)
    stop
    sleep 2
    start
    ;;
  status)
    # This needs improved to use the pid
    status
    exit $?
    ;;
  *)
    echo "Usage: $0 {start|stop|restart}"
    exit 1
    ;;
esac

exit 0
