# templates/teamcity-agent.erb
#! /bin/sh
### BEGIN INIT INFO
# Provides: <%= @service %>
# Required-Start: $local_fs $remote_fs $syslog $network
# Required-Stop: $local_fs $remote_fs $syslog $network
# Default-Start: 2 3 4 5
# Default-Stop: 0 1 6
# Short-Description: TeamCity Build Agent
# Description: TeamCity Build Agent
### END INIT INFO

set -e

PATH=/sbin:/bin:/usr/bin
USER=<%= @user %>
SCRIPT=<%= @home %>/bin/agent.sh
NAME=<%= @agent_name %>
PIDFILE=<%= @temp_dir %>/buildAgent.pid

case "$1" in
  start)
    su $USER --login -c "$SCRIPT start"
    ;;
  stop)
    su $USER --login -c "$SCRIPT stop"
    ;;
  restart)
    su $USER --login -c "$SCRIPT stop"
    sleep 2
    su $USER --login -c "$SCRIPT start"
    ;;
  status)
    # This needs improved to use the pid
    ps -p $(cat $PIDFILE) > /dev/null
    exit $?
    ;;
  *)
    echo "Usage: $0 {start|stop|restart}"
    exit 1
    ;;
esac

exit 0
